# templates/default/storage/etc/udev/rules.d/10_swift.rules.erb
# filter removable devices?

SUBSYSTEM=="block", ACTION=="add", OWNER="swift", \
  RUN+="/usr/local/bin/swift-udev-mount.sh %k"
